# examples/systemd/secrets-sync.service
[Unit]
Description=Docker Secrets Sync Service
Documentation=https://github.com/ohauer/docker-secrets
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/secrets-sync --config /etc/secrets-sync/config.yaml
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=5s

# Environment file (optional)
EnvironmentFile=-/etc/default/secrets-sync

# Security hardening
DynamicUser=yes
StateDirectory=secrets-sync
ConfigurationDirectory=secrets-sync
RuntimeDirectory=secrets-sync

# Filesystem protection
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/var/lib/secrets-sync
PrivateTmp=yes
NoNewPrivileges=yes

# Kernel protection
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectKernelLogs=yes
ProtectControlGroups=yes

# Capabilities
CapabilityBoundingSet=
AmbientCapabilities=

# System calls
SystemCallFilter=@system-service
SystemCallFilter=~@privileged @resources
SystemCallErrorNumber=EPERM

# Network
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
IPAddressDeny=any
IPAddressAllow=localhost
# Uncomment and adjust for your Vault server
#IPAddressAllow=10.0.0.0/8

# Misc
RestrictNamespaces=yes
LockPersonality=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
RemoveIPC=yes
PrivateMounts=yes
ProtectClock=yes

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=secrets-sync

[Install]
WantedBy=multi-user.target
